# server-scripts/neil.crontab
03 2 * * * /home/neil/swintronics/server-scripts/backup-immich.sh
53 3 * * * /home/neil/swintronics/server-scripts/backup-paperless.sh
